# packaging/postinstall.sh
#!/bin/bash
systemctl daemon-reload
systemctl enable vrouter-agent
